# samples/systemd-timers/uablock.timer
[Unit]
Description=uablock auto-updater

[Timer]
OnCalendar=hourly
Unit=uablock.service
Persistent=true

[Install]
WantedBy=multi-user.target